# config/systemd/puma.service.erb
<% require_relative "config/deploy/paths" %>

# DO NOT CHANGE THIS FILE MANUALLY
# IT IS AUTO-GENERATED BY CAPISTRANO FROM TEMPLATE
#
# This file will be installed as code2travel.com_puma.service into /etc/systemd/system/
# Upload it with `cap production systemd:puma:setup`
#
# This file will be installed as <%= $puma_systemd_service %> into /etc/systemd/system/
# Upload it with `cap production systemd:puma:setup`
#
# Puma reference doc:
# https://github.com/puma/puma/tree/master/docs
# https://github.com/puma/puma/blob/master/docs/systemd.md
# https://github.com/puma/puma/blob/master/docs/signals.md#puma-signals
# https://github.com/puma/puma/blob/master/docs/restart.md

[Unit]
Description=Puma application server for <%= $domain %>
After=network.target

<% # Uncomment for socket activation (see below) %>
<% # Requires=puma.socket %>

[Service]
<% # Foreground process (do not use --daemon in ExecStart or config.rb) %>
Type=simple

<% # Preferably configure a non-privileged user %>
User=<%= $deploy_user %>
<% # Group=examplegroup %>

KillMode=mixed
<% # KillMode=process %>

Restart=always
RestartSec=10
<% # Restart=always %>

StartLimitBurst=2

TimeoutSec=10
TimeoutStartSec=5
TimeoutStopSec=10

ExecStart=<%= $puma_command_start %>
ExecStop=<%= $puma_command_stop %>
ExecReload=<%= $puma_command_restart %>

WorkingDirectory=<%= $work_dir %>

PIDFile=<%= $puma_pid_path %>
<% # may conflict with setting inside puma.rb %>

PrivateTmp=false
<% # So that NGINX could access that file otherwise private %>

Environment="PATH=<%= $puma_systemd_bins.join(':') %>"
<%
  # prepending ruby installation dir (like `/opt/rubies/ruby-2.7.1/bin`) is mandatory,
  # otherwise we get
  # > /usr/bin/env: 'ruby': No such file or directory"
  #
  # prepending node bin dir is mandarory,
  # otherwise we get
  # > ! Unable to load application: ExecJS::RuntimeUnavailable:
  # > Could not find a JavaScript runtime. See https://github.com/rails/execjs for a list of available runtimes.
  # > /var/www/sergeypedan.ru/shared/bundle/ruby/2.7.0/gems/execjs-2.7.0/lib/execjs/runtimes.rb:58:in `autodetect':
  # > Could not find a JavaScript runtime.
  # > See https://github.com/rails/execjs for a list of available runtimes. (ExecJS::RuntimeUnavailable)
%>
Environment="PUMA_DEBUG=1"
Environment="PWD=<%= $work_dir %>"
Environment="RAILS_ENV=<%= $rails_env %>"
Environment="RAKE_ENV=<%= $rails_env %>"

<% # EnvironmentFile=#{$work_dir}/.env %>

[Install]
WantedBy=multi-user.target
